# deploy/qq-agent.service
[Unit]
Description=纳西妲 QQ Bot AI Agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300

[Service]
Type=simple
User=REPLACE_WITH_YOUR_USER
Group=REPLACE_WITH_YOUR_GROUP
WorkingDirectory=REPLACE_WITH_PROJECT_PATH

# 可选：检查数据目录是否就绪
# ExecStartPre=/bin/bash -c "test -d /data/db || exit 1"

ExecStart=REPLACE_WITH_PROJECT_PATH/.venv/bin/python qq_bot_adapter.py
Restart=always
RestartSec=10
StartLimitBurst=20

EnvironmentFile=REPLACE_WITH_PROJECT_PATH/.env

# 资源限制
MemoryHigh=1.2G
MemoryMax=1.5G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
